# docker/supervisord-allinone.conf
; ==============================================================================
; SUPERVISORD ALL-IN-ONE PARA CPS EASYPANEL
; 
; Gestiona TODOS los servicios en un solo contenedor:
; - MariaDB (MySQL)
; - Redis
; - Nginx
; - PHP-FPM
; - Laravel Queue Worker
; - Laravel Scheduler
; ==============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

; ==============================================================================
; SERVICIOS DE BASE DE DATOS Y CACHÉ (Prioridad Alta)
; ==============================================================================

[program:mariadb]
command=/usr/bin/mysqld_safe --datadir=/var/lib/mysql
user=mysql
autostart=true
autorestart=true
priority=1
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/mariadb.log
stderr_logfile=/var/log/supervisor/mariadb-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:redis]
command=/usr/bin/redis-server --daemonize no --bind 127.0.0.1 --port 6379 --loglevel warning
user=root
autostart=true
autorestart=true
priority=2
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ==============================================================================
; SERVICIOS WEB (Prioridad Media)
; ==============================================================================

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
user=root
autostart=true
autorestart=true
priority=15
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ==============================================================================
; SERVICIOS DE LARAVEL (Prioridad Baja - dependen de los anteriores)
; ==============================================================================

[program:laravel-queue]
; Ejecutar como root para evitar problemas con ionCube
command=/bin/bash -c "sleep 60 && cd /var/www && php artisan queue:work redis --sleep=3 --tries=3 --timeout=90 --max-jobs=1000 --max-time=3600"
directory=/var/www
user=root
autostart=true
autorestart=true
priority=30
startsecs=65
startretries=5
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/laravel-queue.log
stderr_logfile=/var/log/supervisor/laravel-queue-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/var/www",USER="root"

[program:laravel-scheduler]
; Ejecutar como root para evitar problemas con ionCube
command=/bin/bash -c "sleep 60 && cd /var/www && php artisan schedule:work"
directory=/var/www
user=root
autostart=true
autorestart=true
priority=35
startsecs=65
startretries=5
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stderr_logfile=/var/log/supervisor/laravel-scheduler-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/var/www",USER="root"

; ==============================================================================
; CONFIGURACIÓN DE RPC Y CONTROL
; ==============================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ==============================================================================
; GRUPOS DE SERVICIOS
; ==============================================================================

[group:database]
programs=mariadb,redis
priority=1

[group:web]
programs=nginx,php-fpm
priority=10

[group:laravel]
programs=laravel-queue,laravel-scheduler
priority=30

; ==============================================================================
; NOTAS
; ==============================================================================
; - MariaDB inicia primero (priority=1)
; - Redis inicia segundo (priority=2)
; - Nginx y PHP-FPM inician después (priority=10,15)
; - Laravel queue/scheduler esperan 30s antes de iniciar (priority=30,35)
